# supervisor.conf
[program:foreman]
user=root
command=/opt/start_foreman
redirect_stderr=true
autorestart=true
